# deploy/linux/dynamicqr.service.example
[Unit]
Description=DynamicQR Gunicorn Service
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/dynamicqr
EnvironmentFile=/etc/dynamicqr.env
ExecStart=/opt/dynamicqr/.venv/bin/gunicorn qr_project.wsgi:application --bind 127.0.0.1:8000 --workers 3 --timeout 120
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
